# crontab
*/5 * * * * /opt/minecraft/server/restart.sh
