# crontab
*/1 * * * * /opt/minecraft/server/restart.sh 
